# deploy/systemd-user/oneredpaperclip-horizon.service
[Unit]
Description=One Red Paperclip Horizon
After=network.target
ConditionPathExists=/home/deployer/oneredpaperclip/current

[Service]
Type=simple
WorkingDirectory=/home/deployer/oneredpaperclip/current
ExecStart=/usr/bin/php artisan horizon
Restart=always
RestartSec=3
TimeoutStopSec=90
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
